# systemd/coati-payroll-worker.service
# SPDX-License-Identifier: Apache-2.0
# SPDX-FileCopyrightText: 2025 - 2026 BMO Soluciones, S.A.
#
# Coati Payroll - Dramatiq Worker Service
#
# This systemd unit file manages the Dramatiq background worker for processing
# payroll calculations asynchronously. This service is REQUIRED when QUEUE_ENABLED=1
# and background processing is configured.
#
# Installation:
#   1. Copy this file to /etc/systemd/system/coati-payroll-worker.service
#   2. Edit the Environment variables below to match your deployment (must match main service)
#   3. Run: sudo systemctl daemon-reload
#   4. Run: sudo systemctl enable coati-payroll-worker
#   5. Run: sudo systemctl start coati-payroll-worker
#
# Usage:
#   - Start:   sudo systemctl start coati-payroll-worker
#   - Stop:    sudo systemctl stop coati-payroll-worker
#   - Restart: sudo systemctl restart coati-payroll-worker
#   - Status:  sudo systemctl status coati-payroll-worker
#   - Logs:    sudo journalctl -u coati-payroll-worker -f
#
# Note: This service should be started BEFORE the main application to ensure
#       background jobs can be processed immediately.

[Unit]
Description=Coati Payroll Dramatiq Worker
Documentation=https://github.com/williamjmorenor/coati-payroll
After=network.target redis.service
Before=coati-payroll.service
Requires=redis.service

[Service]
Type=simple
User=coati
Group=coati
WorkingDirectory=/home/coati/coati-payroll

# Path to Python virtual environment
Environment="PATH=/home/coati/coati-payroll/venv/bin"

# Database configuration (same as main service)
Environment="DATABASE_URL=postgresql://coati_user:CHANGE_ME@localhost:5432/coati_db"

# Redis configuration (REQUIRED for Dramatiq)
Environment="REDIS_URL=redis://localhost:6379/0"

# Queue system configuration
Environment="QUEUE_ENABLED=1"

# Worker configuration (adjust based on server capacity)
Environment="DRAMATIQ_WORKER_THREADS=8"
Environment="DRAMATIQ_WORKER_PROCESSES=2"

# Start the Dramatiq worker
ExecStart=/home/coati/coati-payroll/venv/bin/dramatiq coati_payroll.queue.tasks --threads 8 --processes 2

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coati-payroll-worker

[Install]
WantedBy=multi-user.target
